# etc/supervisor/worker.ini
[program:worker]
command=php /app/yii queue/listen
priority=1
numprocs=5
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
